# hidindicator-1.0.0/rules/35-get-usb-device-info.rules
ATTR{bDeviceClass}=="00", ACTION=="add" , RUN+="/usr/bin/store '%s{manufacturer}' '%s{product}"
