# hidindicator-1.0.0/rules/90-inform-user-on-hid.rules
ATTR{bInterfaceClass}=="03", ACTION=="add" , RUN+="/usr/bin/notify '%s{manufacturer}'"


